# docker/supervisord.conf
[supervisord]
nodaemon                = true
logfile                  = /tmp/supervisord.log
pidfile                  = /tmp/supervisord.pid
logfile_maxbytes         = 5MB
stdout_logfile           = /dev/stdout
stdout_logfile_maxbytes  = 0
stderr_logfile           = /dev/stderr
stderr_logfile_maxbytes  = 0
stdout_events_enabled   = true
stderr_events_enabled   = true

[program:nginx]
command                 = /usr/sbin/nginx
autorestart             = true
stdout_logfile           = /dev/stdout
stdout_logfile_maxbytes  = 0
stderr_logfile           = /dev/stderr
stderr_logfile_maxbytes  = 0
stdout_events_enabled   = true
stderr_events_enabled   = true

[program:php-fpm7]
command                 = php-fpm7 -F -c /etc/php7/php.ini -y /etc/php7/php-fpm.conf
autorestart             = true
stdout_logfile           = /dev/stdout
stdout_logfile_maxbytes  = 0
stderr_logfile           = /dev/stderr
stderr_logfile_maxbytes  = 0
stdout_events_enabled   = true
stderr_events_enabled   = true